# irexec.conf
[program:irexec]
command=irexec -D info /etc/lirc/irexec.lircrc
autostart=false
startsecs=0
dependent_startup=true
dependent_startup_wait_for=lircd:running
