# lircd.conf
[program:lircd]
command=lircd
autostart=false
startsecs=3
dependent_startup=true
